# luci-app-adguardhome/root/etc/hotplug.d/iface/99-AdGuardHome
#!/bin/sh

[ "$ACTION" = ifup ] && /etc/init.d/AdGuardHome restart
